# crontab
* * * * *       cd /root/prj/devcron && ./periodic.lua
23 59 * * *     /root/prj/devcron/logrotate
