# crontab
* * * * *       cd /root/prj/devcron && ./periodic.lua
0 0 * * *       mv /tmp/periodic.log /tmp/periodic.log.bak

